# services/squid3-noauth.conf
description "Squid3 proxy without auth"
start on filesystem and started docker
respawn
script
 docker start -a squid-noauth
end script


pre-stop script
 docker stop squid-noauth
end script
